# backend/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:laravel-server]
command=php artisan serve --host=0.0.0.0 --port=8000
directory=/var/www
user=www-data
autostart=true
autorestart=true
stdout_logfile=/var/www/storage/logs/laravel-server.log
stderr_logfile=/var/www/storage/logs/laravel-server-error.log

[program:laravel-scheduler]
command=/bin/bash -c "while true; do php artisan schedule:run >> /var/www/storage/logs/scheduler.log 2>&1; sleep 60; done"
directory=/var/www
user=www-data
autostart=true
autorestart=true
stdout_logfile=/var/www/storage/logs/scheduler-supervisor.log
stderr_logfile=/var/www/storage/logs/scheduler-error.log
